# support/supervisord/cert_manage.ini
[program:cert_manage]
directory=/opt/cert_manage
environment=PATH="/opt/py3/bin:%(ENV_PATH)s"
command=/opt/py3/bin/python run_server.py all
pidfile=/opt/cert_manage/cert_manage.pid
autostart=true
autorestart=true
stdout_logfile=/opt/cert_manage/logs/cert_manage.log
stderr_logfile=/opt/cert_manage/logs/cert_manage_error.log